# rootfs/etc/s6-overlay/s6-rc.d/telegraf_socat_978/run
#!/bin/sh
exec /scripts/telegraf_socat_978
